# .config/zsh/profile.d/90-home-bin-path.zsh
# Add ~/bin to beginning of stack (mosty last) to allow
# overriding any system-installed executables
path=("$HOME/bin" $path)
